# lede-imagebuilder-ar71xx-generic.Linux-x86_64/build_dir/target-mips_24kc_musl-1.1.16/root-ar71xx/usr/lib/opkg/info/ChinaDNS.list
/etc/init.d/chinadns
/usr/bin/chinadns
/etc/config/chinadns
/etc/chinadns_chnroute.txt
