# lede-imagebuilder-ar71xx-generic.Linux-x86_64/build_dir/target-mips_24kc_musl-1.1.16/root-ar71xx/usr/lib/opkg/info/luci-app-shadowsocks.list
/etc/config/shadowsocks
/usr/bin/ss-rules
/usr/lib/lua/luci/controller/shadowsocks.lua
/usr/lib/lua/luci/model/cbi/shadowsocks/general.lua
/usr/lib/lua/luci/view/shadowsocks/general.htm
/usr/lib/lua/luci/view/shadowsocks/dynamiclist.htm
/etc/init.d/shadowsocks
/usr/lib/lua/luci/i18n/shadowsocks.zh-cn.lmo
/usr/lib/lua/luci/model/cbi/shadowsocks/servers-details.lua
/usr/lib/lua/luci/model/cbi/shadowsocks/servers.lua
/etc/uci-defaults/luci-shadowsocks
/usr/lib/lua/luci/model/cbi/shadowsocks/access-control.lua
/usr/lib/lua/luci/view/shadowsocks/servers-details.htm
